# rsyslog-source/files/supervisor.d/rsyslogd.conf
[program:rsyslogd]
priority=1
command=/usr/sbin/rsyslogd -n
